# ansible/roles/truenas_snapshots/templates/sanoid.timer.j2
[Unit]
Description=Run Sanoid every 15 minutes
Requires=sanoid.service

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target